# templates/container-config.path
[Unit]
Description=regenerate container if file changes
ConditionPathExists=!/run/containers/{{ deploy_app.name }}/block

[Path]
PathChanged=/etc/containers/config/{{ deploy_app.name }}
PathExistsGlob=/etc/containers/config/{{ deploy_app.name }}/*
Unit={{ deploy_app.name }}--{{ deploy_app.host_port }}.service

[Install]
WantedBy=app-run.target
